# .config/systemd/user/check_battery.service
[Unit]
Description=Check Battery Status

[Service]
ExecStart=%h/.config/scripts/check_battery.sh

[Install]
WantedBy=default.target
