# .config/systemd/user/check_battery.timer
[Unit]
Description=Check Battery Status Timer

[Timer]
OnActiveSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
